# infra/systemd/bumpabestie-disk-usage.service
[Unit]
Description=Check Bumpa Bestie host disk and inode usage
After=local-fs.target

[Service]
Type=oneshot
User=bumpabestie
Group=bumpabestie
WorkingDirectory=/opt/bumpabestie
Environment=BUMPABESTIE_DISK_THRESHOLD_PERCENT=85
Environment=BUMPABESTIE_DISK_PATHS=/
Environment=BUMPABESTIE_ALERT_HOOK=/opt/bumpabestie/scripts/send_ops_alert.py
Environment=BUMPABESTIE_ALERT_CONFIG_FILE=/etc/bumpabestie/alerts.json
ExecStart=/usr/bin/python3 /opt/bumpabestie/scripts/check_disk_usage.py
TimeoutStartSec=45
UMask=0077
StandardOutput=journal
StandardError=journal

# /etc/bumpabestie/alerts.json contains only the fixed endpoint, secret-file
# path and bounded retry settings. The secret value remains in its 0400 file.
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
RemoveIPC=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
