# infra/systemd/bumpabestie-backup.service
[Unit]
Description=Bumpa Bestie local backup and optional off-host handoff
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
User=bumpabestie
WorkingDirectory=/opt/bumpabestie
ExecStart=/opt/bumpabestie/scripts/scheduled_backup.sh
ExecStartPost=/opt/bumpabestie/scripts/offsite_backup.sh --env-file /opt/bumpabestie/.env.production
Nice=10
