# infra/systemd/bumpabestie-disk-usage.timer
[Unit]
Description=Check Bumpa Bestie host disk usage every five minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=30s
RandomizedDelaySec=30s
Persistent=true
Unit=bumpabestie-disk-usage.service

[Install]
WantedBy=timers.target
