# infra/systemd/bumpabestie-backup.timer
[Unit]
Description=Run Bumpa Bestie backup nightly

[Timer]
OnCalendar=*-*-* 02:30:00 UTC
RandomizedDelaySec=900
Persistent=true
Unit=bumpabestie-backup.service

[Install]
WantedBy=timers.target
