# infra/systemd/docker.service.d/10-bumpabestie-cloudflare-origin-pregate.conf
[Unit]
Requires=bumpabestie-cloudflare-origin-pregate.service
After=bumpabestie-cloudflare-origin-pregate.service
